# files/c-icap.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 OpenWrt.org

START=90
STOP=10

USE_PROCD=1
PROG=/usr/bin/c-icap
CONFIG_FILE=/etc/c-icap/c-icap.conf
PID_FILE=/var/run/c-icap/c-icap.pid

start_service() {
        "${PROG}" -f "$CONFIG_FILE"
}

stop_service() {
	[ ! -f "$PID_FILE" ] && { echo "c-icap is not running"; exit 1; }
	PID=$(cat "$PID_FILE")
        kill ${PID}
}

